# cafeext/skills/restart_gateway/scripts/com.nanobot.gateway.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanobot.gateway.watchdog</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/leolin/Desktop/Work/nanobot__ava/nanobot/skills/restart_gateway/scripts/gateway_watchdog.sh</string>
        <string>--silent</string>
    </array>
    
    <!-- 每 5 分钟执行一次 -->
    <key>StartInterval</key>
    <integer>300</integer>
    
    <!-- 系统启动后也执行一次 -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- 环境变量 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/leolin</string>
    </dict>
    
    <!-- 日志输出 -->
    <key>StandardOutPath</key>
    <string>/tmp/gateway_watchdog_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/gateway_watchdog_stderr.log</string>
    
    <!-- 低优先级 -->
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
